# init-files/backup.cron
0 5 * * 1-5 root /bin/bash /var/jenkins_home/backup/backup.sh